# docker/php-fpm.conf
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log
